# ansible/roles/step_ca/files/step-ca.service
[Unit]
Description=step-ca Certificate Authority
After=network.target

[Service]
Type=simple
User=stepca
Group=stepca
ExecStart=/usr/bin/step-ca /etc/step-ca/config/ca.json --password-file /etc/step-ca/password.txt
Restart=on-failure

[Install]
WantedBy=multi-user.target
